# configs/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
[program:httpd]
command=/etc/apache2/foreground.sh
stopsignal=6
stdout_logfile=/var/log/supervisor/%(program_name)s.log
sterr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
